# skynet/dotfiles/.config/systemd/user/udiskie.service
[Unit]
Description=udiskie - automounter for removable media

[Service]
Type=simple
ExecStart=/usr/bin/udiskie

[Install]
WantedBy=default.target
